# etc/profile.d/neovide.sh
#!/bin/sh
export NEOVIDE_MULTIGRID=true
export NEOVIDE_FRAMELESS=true
